# systemd/module-registry.service.d/10-params.conf
# Load the supervisor's managed_nodes list (installed at
# /etc/systemd/system/module-registry.service.d/10-params.conf).
# Params file lives at /etc/friday/core_hub.params.yaml (source: config/core_hub.params.yaml).
[Service]
ExecStart=
ExecStart=/opt/friday/bin/ros2-run friday_core_hub core_hub --ros-args -p authority_holder:=MARK1-CORE-001 --params-file /etc/friday/core_hub.params.yaml
