# systemd/module-registry.service.d/20-udpv4.conf
[Service]
Environment=FASTDDS_BUILTIN_TRANSPORTS=UDPv4
