# systemd/mosquitto-internal.service
[Unit]
Description=Mosquitto MQTT broker — internal bridge (Core ↔ Telemetry Gateway)
Documentation=https://github.com/Friday-Labs-Inc/friday-core-os
After=network-online.target
Wants=network-online.target
PartOf=friday-core-os.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/sbin/mosquitto -c /etc/friday/mosquitto-internal.conf
Restart=on-failure
RestartSec=3
User=mosquitto
Group=mosquitto

[Install]
WantedBy=friday-core-os.target
