# systemd/phone-bridge.service
[Unit]
Description=Friday Labs phone pod bridge (GPS + IMU over WiFi)
PartOf=friday-core-os.target
After=module-registry.service gpsd.service network-online.target
[Service]
Environment=RMW_IMPLEMENTATION=rmw_fastrtps_cpp
Environment=ROS_LOCALHOST_ONLY=0
Environment=ROS_DOMAIN_ID=42
Environment=ROS_HOME=/var/lib/friday/ros
Environment=FASTDDS_BUILTIN_TRANSPORTS=UDPv4
ExecStart=/opt/friday/bin/ros2-run friday_phone_bridge phone_bridge
Restart=always
RestartSec=3
[Install]
WantedBy=friday-core-os.target
